# systemd/coding-partner.service
[Unit]
Description=Coding Partner - Feishu Vibe Coding Bot
After=network.target

[Service]
Type=simple
User=yuanzhihong.chiwei
WorkingDirectory=/home/yuanzhihong.chiwei/projects/coding-partner
ExecStart=/bin/bash -lc 'exec /home/yuanzhihong.chiwei/.local/bin/uv run python -m coding_partner.main'
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/yuanzhihong.chiwei/projects/coding-partner/.env

[Install]
WantedBy=multi-user.target
